# controller/99-newtonwars.rules
ACTION=="add", KERNEL=="ttyACM*", \
	RUN+="/home/pi/controller.sh"
